# Backup - nG1 Servers/full_bk.timer
[Unit]
Description="Run full_bk.service" 

[Timer]
OnBootSec=5min
OnUnitActiveSec=4hr
Unit=full_bk.service

[Install]
WantedBy=multi-user.target
